# .config/systemd/user/eww-daemon.service
[Unit]
Description=Eww daemon

[Service]
Type=forking
ExecStart=/usr/bin/eww daemon
ExecStop=/usr/bin/eww kill

[Install]
WantedBy=multi-user.target
